# misc/uwsgi.conf
# file: /etc/init/uwsgi.conf
description "uWSGI starter"

start on (local-filesystems and runlevel [2345])
stop on runlevel [016]

respawn

# home - is the path to our virtualenv directory
# pythonpath - the path to our django application
# module - the wsgi handler python script

export IS_PRODUCTION=true
BASEDIR=/home/ubuntu/sotmjp-website

exec ${BASEDIR}/env/sotmjp/bin/uwsgi --chdir=${BASEDIR} \
  -s 127.0.0.1:3031  \
  --module=symposion.wsgi:application \
  --home ${BASEDIR}/env/sotmjp \
  --master \
  --pidfile=/tmp/uwsgi-master.pid \
  --processes=2 \
  --daemonize=/home/ubuntu/uwsgi.log
